# templates/default/clojure_web.upstart.conf.erb
description "An Uberjarred web application"

start on [2345]
stop on [!12345]

respawn
respawn limit 2 5

kill timeout 300

setuid <%= @user %>
setgid <%= @group %>

exec java -Xms<%= @init_heap %> -Xmx<%= @max_heap %> -jar /opt/<%= @name %>.jar /opt/<%= @name %>.edn

post-start script
  PID=`status <%= @name %> | egrep -oi '([0-9]+)$' | head -n1`
  echo $PID > /opt/<%= @name %>/<%= @name %>.pid
end script

post-stop script
  rm -f /opt/<%= @name %>/<%= @name %>.pid
end script
